# monitoring-01/Add/crontab
* * * * * python3 /home/atimofeev/Scripts/logs.py